# yedekler/ExcPanel.TransferAgent_2026-06-22_12-16-29/deploy/linux/excpanel-transfer-agent.service
[Unit]
Description=ExcPanel Transfer Agent
After=network.target

[Service]
Type=simple
User=excpanel-agent
Group=excpanel-agent
WorkingDirectory=/opt/excpanel-transfer-agent
ExecStart=/opt/excpanel-transfer-agent/ExcPanel.TransferAgent
Restart=on-failure
RestartSec=5
# Agent must invoke the privileged helper via sudo; NoNewPrivileges blocks elevation.
NoNewPrivileges=false
PrivateTmp=true
# Helper performs disk, fstab and mount operations outside the install directory.
ProtectSystem=false
ProtectHome=false
MountPropagation=shared

[Install]
WantedBy=multi-user.target
